# nomic-embed@.service
[Unit]
Description=Nomic Embed pool instance on port %i (GPU bulk ingest)
After=network.target

[Service]
Type=simple
User=kevyn
Restart=always
RestartSec=3
EnvironmentFile=-/opt/ai/config/nomic-embed.env
EnvironmentFile=-/opt/ai/config/nomic-embed-scale.env
Environment=LLAMA_SERVER_BIN=/opt/ai/bin/llama-server
Environment=NOMIC_EMBED_MODEL=/opt/ai/models/embed/nomic-embed-text-v1.5.Q8_0.gguf
Environment=NOMIC_GPU_LAYERS=99
Environment=NOMIC_CONTEXT=8096
Environment=NOMIC_POOL_PARALLEL=16

ExecStart=/opt/ai/bin/llama-server \
    -m $NOMIC_EMBED_MODEL \
    --embedding \
    --pooling mean \
    -ngl $NOMIC_GPU_LAYERS \
    -c $NOMIC_CONTEXT \
    --parallel $NOMIC_POOL_PARALLEL \
    -b 2048 \
    -ub 2048 \
    --host 127.0.0.1 \
    --port %i \
    --log-disable

[Install]
WantedBy=multi-user.target
