# nomic-embed-scale.service
[Unit]
Description=Scale nomic-embed GPU pool to fit free VRAM
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Root so --apply can write /opt/ai/config/nomic-embed-pool.env and manage units.
User=root
WorkingDirectory=/opt/ai/repo/rag_proxy
EnvironmentFile=-/opt/ai/config/nomic-embed-scale.env
EnvironmentFile=-/opt/ai/config/rag-proxy.env
EnvironmentFile=-/opt/ai/config/rag-admin.env
ExecStart=/opt/ai/venv/bin/python scripts/scale_ingest_capacity.py --apply
RemainAfterExit=yes

# Edit WorkingDirectory and ExecStart if repo lives elsewhere (e.g. /home/kevyn/rag_proxy).

[Install]
WantedBy=multi-user.target
