# rag-proxy.service
[Unit]
Description=RAG proxy for llama-swap (Qdrant context injection)
After=network.target
# nomic-embed starts on demand (EMBED_ON_DEMAND) when RAG or ingest needs it
Before=llama-swap.service

[Service]
Type=simple
User=kevyn
WorkingDirectory=/home/kevyn/rag_proxy
EnvironmentFile=-/home/kevyn/rag_proxy/.env
ExecStart=/home/kevyn/rag_proxy/.venv/bin/python rag_proxy.py
Restart=always
RestartSec=3
TimeoutStopSec=30
KillMode=mixed

[Install]
WantedBy=multi-user.target
