# rag-proxy.service
[Unit]
Description=RAG proxy for llama-swap (Qdrant context injection)
After=network.target
# nomic-embed starts on demand (EMBED_ON_DEMAND) when RAG or ingest needs it
Before=llama-swap.service

[Service]
Type=simple
User=kevyn
WorkingDirectory=/home/kevyn/rag-proxy
EnvironmentFile=-/opt/ai/config/rag-proxy.env
ExecStart=/opt/ai/venv/bin/python rag_proxy.py
Restart=always
RestartSec=3
TimeoutStopSec=30
KillMode=mixed

# Do not cp blindly — run: bash scripts/install-systemd-units.sh

[Install]
WantedBy=multi-user.target
