# rag-admin.service
[Unit]
Description=RAG admin UI and dense ingest worker
After=network.target

[Service]
Type=simple
User=kevyn
WorkingDirectory=/home/kevyn/rag-proxy
EnvironmentFile=-/opt/ai/config/rag-admin.env
ExecStart=/opt/ai/venv/bin/python -m rag_admin
Restart=always
RestartSec=3
TimeoutStopSec=30
KillMode=mixed

# Do not cp blindly — run: bash scripts/install-systemd-units.sh

[Install]
WantedBy=multi-user.target
