# rag-proxy.service
[Unit]
Description=RAG proxy for llama-swap (Qdrant context injection)
After=network.target nomic-embed.service
Wants=nomic-embed.service
# Start before clients use it; llama-swap can start in parallel
Before=llama-swap.service

[Service]
Type=simple
User=kevyn
WorkingDirectory=/home/kevyn/rag_proxy
EnvironmentFile=-/home/kevyn/rag_proxy/.env
ExecStart=/home/kevyn/rag_proxy/.venv/bin/python rag_proxy.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
